# crontab
*/1 * * * * /Users/student/.nvm/versions/node/v6.9.5/bin/node /Users/student/Desktop/hrsf86-web-historian/workers/htmlfetcher.js